# board/common/rootfs-overlay/usr/lib/systemd/system/rngd.service.d/override.conf
[Unit]
DefaultDependencies=no
After=systemd-udev-settle.service
Before=sysinit.target shutdown.target
Conflicts=shutdown.target
