# board/common/rootfs-overlay/usr/lib/systemd/system/boot.automount
[Automount]
Where=/boot
DirectoryMode=0700
TimeoutIdleSec=5

[Install]
WantedBy=local-fs.target
